# systemd/.config/systemd/user/emacsd.service
[Unit]
Description=Emacs Deamon Job

[Service]
Type=simple
ExecStart=/usr/bin/emacs --debug-init --fg-daemon
Restart=always

